# ops/prod/systemd/parallel-reminder-email-queue.service
[Unit]
Description=parallel reminder email queue
After=network.target

[Service]
Environment=PORT=4000
Environment=NODE_ENV=production
WorkingDirectory=/home/ec2-user/main/server
Type=simple
User=ec2-user
ExecStart=/usr/bin/node /home/ec2-user/main/server/dist/workers/reminder-email.js
Restart=always

[Install]
WantedBy=multi-user.target
